# deltamike-node.conf
start on startup

respawn
respawn limit 99 5

script
  cd /home/ubuntu/apps/deltamike_node/current/
  exec npm run start-server-prod > /dev/null 2>&1 &
end script
